# debbuilder/nubogateway2/DEBIAN/postrm
#!/bin/sh

if [ "$1" = "remove" ]; then
	systemctl disable nubogateway ||:
fi

#Restart after every install/update
service rsyslog reload > /dev/null 2>&1 ||:
