# docker_build/etc/supervisor/conf.d/guacd.conf
[program:guacd]
directory=/usr/local/guacamole
command=/usr/local/guacamole/sbin/guacd -b 0.0.0.0 -L info -f
user=guacd
autorestart=true
priority=400
