# docker_build/etc/supervisor/conf.d/nubogateway.conf
[program:nubogateway]
directory=/opt/nubogateway
command=/usr/bin/node dist/gateway.js
user=root
autorestart=true
priority=400
