# docker_build/etc/supervisor/conf.d/nubogateway.conf
[program:nubogateway]
directory=/opt/nubogateway
command=/usr/bin/node dist/gateway.js
user=nubo
environment=HOME="/home/nubo",USER="nubo"
autorestart=true
priority=400
